# docker/supervisord/conf/crontab.conf
### to schedule run, uncomment all
[program:crontab]
command=/bin/bash init.sh
directory=/etc/supervisord.d/
autostart=true
autorestart=false
startsecs = 0
startretries = 1



